# crontab
*/10 * * * * /usr/bin/fetchmail > /dev/null 2>&1